# ansible/templates/docsend.service.j2
[Unit]
Description=Docsend Application
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User={{ app_user }}
Group={{ app_user }}
WorkingDirectory={{ app_dir }}
Environment="DOCSEND_CONFIG={{ config_dir }}/config.toml"
Environment="RUST_LOG={{ lookup('env', 'RUST_LOG') | default('info', true) }}"
ExecStart={{ app_dir }}/bin/docsend
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=docsend

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths={{ data_dir }} {{ log_dir }}

[Install]
WantedBy=multi-user.target
